# utils/service/scheduler/run
#!/bin/bash

cd ../../../
node lib/awesomud.js
